# crontab
# IX-Tips scheduled jobs — run by supercronic (replaces Celery beat).
# Times are interpreted in the TZ env var set on the machine (Africa/Nairobi).
#
# min hour dom mon dow   command
  30   5    *   *   *     cd /code && python manage.py run_task training
  45   5    *   *   *     cd /code && python manage.py run_task warmform
  0    6    *   *   *     cd /code && python manage.py run_task predictions
  20   6    *   *   *     cd /code && python manage.py run_task odds
  30   6    *   *   *     cd /code && python manage.py run_task combo
  1    *    *   *   *     cd /code && python manage.py run_task tables
  0    *    *   *   *     cd /code && python manage.py run_task metadata
  */5  *    *   *   *     cd /code && python manage.py run_task live
